# 2023/projects/02_ProtoFabStats/rapsberrypi_app/data_analysis_services/hour_stats_service/hour_stats.timer
[Unit]
Description=PortoFabStats Day Of Week Stats Timer

[Timer]
OnUnitActiveSec=1h
OnBootSec=30s
AccuracySec=2s

[Install]
WantedBy=timers.target