# 2023/projects/02_ProtoFabStats/rapsberrypi_app/data_analysis_services/data_preprocessing_service/proto_fab_stats_preprocessing.timer
[Unit]
Description=PortoFabStats Preprocessing Timer

[Timer]
OnUnitActiveSec=10m
OnBootSec=30s
AccuracySec=2s

[Install]
WantedBy=timers.target