# 2023/projects/02_ProtoFabStats/rapsberrypi_app/data_analysis_services/day_of_week_stats_service/day_of_week_stats.timer
[Unit]
Description=PortoFabStats Day Of Week Stats Timer

[Timer]
Daily at 4:00am
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target