# 2023/topics/02_rasPI/data_recorder/data_recorder.timer
[Unit]
Description=data recorder

[Timer]
OnUnitActiveSec=30s
OnBootSec=30s
AccuracySec=2s

[Install]
WantedBy=timers.target
